# admin/dev/etc/xinetd.d/vista-rpcbroker
service vista-rpcbroker
{
  port = 9430
  socket_type = stream
  protocol = tcp
  type = UNLISTED
  user = osehra
  server = /home/osehra/bin/rpcbroker.sh
  wait = no
  disable = yes
  per_source = UNLIMITED
  instances = UNLIMITED
}
